# openclaw/systemd/openclaw-auto-backtest-daily.timer
[Unit]
Description=Run OpenClaw Auto Backtest Daily

[Timer]
OnCalendar=Mon..Fri 22:05
Persistent=true
Unit=openclaw-auto-backtest-daily.service

[Install]
WantedBy=timers.target

