# _archive/launchd/com.airivo.openclaw.data-optimize.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.airivo.openclaw.data-optimize</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mac/2026Qlin/scripts/openclaw_data_optimize.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mac/2026Qlin</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/mac/2026Qlin/logs/openclaw/data_optimize.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mac/2026Qlin/logs/openclaw/data_optimize.launchd.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
