# deploy_stock/openclaw-db-sync.timer
[Unit]
Description=Daily stock DB sync from A server

[Timer]
OnCalendar=*-*-* 07:10:00
Persistent=true
Unit=openclaw-db-sync.service

[Install]
WantedBy=timers.target

