# stock_ultimate_system/deploy/aliyun/stock-ultimate-dashboard.service.d/canonical-artifacts.conf
[Service]
Environment=STOCK_ULTIMATE_ARTIFACTS_DIR=/opt/stock-ultimate/app/artifacts
Environment=STOCK_ULTIMATE_EXPERIMENTS_DIR=/opt/stock-ultimate/app/data/experiments
