# stock_ultimate_system/deploy/aliyun/airivo-apex-full-market-candidates.service
[Unit]
Description=Airivo Apex Full-Market Candidate Refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/airivo-apex/app
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/airivo-apex/app
TimeoutStartSec=12h
ExecStart=/opt/airivo-apex/.venv/bin/python /opt/airivo-apex/app/run_update_database.py --config-dir config/server --post-top-candidates --no-post-daily-research --post-universe-size 0 --allow-full-market-post-candidates --post-top-n 20
StandardOutput=append:/var/log/airivo-apex/full_market_candidates.log
StandardError=append:/var/log/airivo-apex/full_market_candidates.err.log
