# stock_ultimate_system/deploy/aliyun/airivo-apex-primary-result-daily-closure.timer
[Unit]
Description=Run Airivo Apex Primary Result Daily Closure at 18:40 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 18:40:00 Asia/Shanghai
Persistent=true
Unit=airivo-apex-primary-result-daily-closure.service

[Install]
WantedBy=timers.target
