# stock_ultimate_system/deploy/aliyun/airivo-apex-t12.service
[Unit]
Description=Airivo Apex T12 Readonly Governance
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/airivo-apex/app
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/airivo-apex/app
ExecStart=/opt/airivo-apex/.venv/bin/python /opt/airivo-apex/app/run_dashboard.py --host 127.0.0.1 --port 18766 --base-path /apex/T12
Restart=always
RestartSec=5
StandardOutput=append:/var/log/airivo-apex/t12.log
StandardError=append:/var/log/airivo-apex/t12.err.log

[Install]
WantedBy=multi-user.target
