# stock_ultimate_system/deploy/aliyun/airivo-apex-primary-result-feedback-loop.timer
[Unit]
Description=Run Airivo Apex Primary Result Feedback Loop

[Timer]
OnCalendar=*-*-* 18:45:00 Asia/Shanghai
Persistent=true
Unit=airivo-apex-primary-result-feedback-loop.service

[Install]
WantedBy=timers.target
