# stock_ultimate_system/deploy/aliyun/airivo-apex-performance-evidence.timer
[Unit]
Description=Run Airivo Apex Performance Evidence Report

[Timer]
OnCalendar=*-*-* 18:58:00 Asia/Shanghai
Persistent=true
Unit=airivo-apex-performance-evidence.service

[Install]
WantedBy=timers.target
