# stock_ultimate_system/deploy/aliyun/stock-ultimate-update.timer
[Unit]
Description=Run Stock Ultimate Daily Update at 17:30 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 17:30:00 Asia/Shanghai
Persistent=true
Unit=stock-ultimate-update.service

[Install]
WantedBy=timers.target
