# stock_ultimate_system/deploy/aliyun/stock-ultimate-t12.service
[Unit]
Description=Stock Ultimate T12 Readonly Governance
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/stock-ultimate/app
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/stock-ultimate/app
ExecStart=/opt/stock-ultimate/.venv/bin/python /opt/stock-ultimate/app/run_dashboard.py --host 127.0.0.1 --port 8766 --base-path /T12
Restart=always
RestartSec=5
StandardOutput=append:/var/log/stock-ultimate/t12.log
StandardError=append:/var/log/stock-ultimate/t12.err.log

[Install]
WantedBy=multi-user.target
