# stock_ultimate_system/deploy/aliyun/stock-ultimate-entry-guard.timer
[Unit]
Description=Refresh Stock Ultimate Entry Guard

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=stock-ultimate-entry-guard.service
Persistent=true

[Install]
WantedBy=timers.target
