# stock_ultimate_system/deploy/aliyun/stock-ultimate-healthcheck.service
[Unit]
Description=Stock Ultimate Automation Health Check
After=network-online.target stock-ultimate-update.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/stock-ultimate/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash /opt/stock-ultimate/app/tools/check_automation_health.sh /opt/stock-ultimate/app/data/experiments/update_status_latest.json
StandardOutput=append:/var/log/stock-ultimate/automation_health.log
StandardError=append:/var/log/stock-ultimate/automation_health.err.log
