# _archive/openclaw/systemd/openclaw-tunnel.service
[Unit]
Description=OpenClaw Tunnel
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/ssh -i /root/.ssh/id_openclaw_tunnel -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -N -L 127.0.0.1:15000:127.0.0.1:5000 root@47.90.160.87
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
