# _archive/openclaw/systemd/openclaw-auto-evolve-opt.timer
[Unit]
Description=Run OpenClaw Auto Evolve Optimize Nightly

[Timer]
OnCalendar=Mon..Fri 22:20
Persistent=true
Unit=openclaw-auto-evolve-opt.service

[Install]
WantedBy=timers.target
